# scripts/hermes_mon/com.yangtianbao.hermes-mon.aggregator.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.yangtianbao.hermes-mon.aggregator</string>

    <key>ProgramArguments</key>
    <array>
        <string>/usr/local/bin/python3</string>
        <string>/Users/yangtb/.hermes/hermes-agent/scripts/hermes_mon/aggregator.py</string>
    </array>

    <!-- 每小时整点运行（第 5 分钟，避免与采集器的整点文件切换竞态） -->
    <key>StartCalendarInterval</key>
    <dict>
        <key>Minute</key>
        <integer>5</integer>
    </dict>

    <key>RunAtLoad</key>
    <false/>

    <key>LowPriorityIO</key>
    <true/>
    <key>Nice</key>
    <integer>10</integer>

    <key>StandardOutPath</key>
    <string>/tmp/hermes-mon-aggregator.log</string>
    <key>StandardErrorPath</key>
    <string>/tmp/hermes-mon-aggregator.log</string>
</dict>
</plist>
